# ansible/playbooks/files/service-files/mount-all.service
[Unit]
Description=Mount all filesystems from fstab
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/mount -a
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
