# ansible/playbooks/files/service-files/qbittorrent-nox.service
[Unit]
Description=qBittorrent-nox Daemon Service
After=network.target

[Service]
User=root
ExecStart=/usr/bin/qbittorrent-nox
Restart=on-failure

[Install]
WantedBy=multi-user.target